# 20_scripts/freebox-vaults-watch.path
[Unit]
Description=Watch ~/Vaults/ for new subdirectories and start Claude sessions
Documentation=file:///home/frimann/Programming/freeBox/10_docs/freebox-services.md

[Path]
PathChanged=%h/Vaults
MakeDirectory=yes

[Install]
WantedBy=default.target
